# etc/systemd/system/uefiboot-update.service
[Unit]
Description=Update UEFI boot options according to installed kernel versions

[Service]
Type=oneshot
ExecStart=/usr/bin/uefiboot-update
