# etc/systemd/system/uefiboot-update.path
[Unit]
Description=Update UEFI boot options according to installed kernel versions

[Path]
PathModified=/vmlinuz
PathModified=/vmlinuz.old

[Install]
WantedBy=multi-user.target